# docker/supervisord.conf
# Kamto — produkční process manager. Magic Containers spustí jeden kontejner, ale appka potřebuje
# dva procesy (php-fpm + nginx). supervisord je řídí a restartuje při pádu; logy jdou na stdout/stderr
# (Bunny sbírá container logy). Běží jen v produkčním image (dev používá samostatné compose služby).
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/run/supervisord.pid

[program:php-fpm]
command=php-fpm -F
# autorestart: když fpm spadne, supervisord ho restartuje NA MÍSTĚ (kontejner běží dál). Pro
# single-user appku je to lepší než cold restart celého kontejneru — nginx pár requestů vrátí 502,
# než fpm naběhne. HEALTHCHECK (Dockerfile) hlásí ne-živý stav, kdyby fpm skončil ve FATAL.
autorestart=true
startretries=3
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
# QUIT = graceful shutdown (dokončí rozdělané requesty) místo defaultního TERM (fast).
stopsignal=QUIT
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
